# supervisor.conf
; Supervisor configuration for Behome queue workers
; Copy this to /etc/supervisor/conf.d/behome.conf on your production server
; Then run: sudo supervisorctl reread && sudo supervisorctl update && sudo supervisorctl start all

[program:behome-worker-emails]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/behome/artisan queue:work database --queue=emails --sleep=3 --tries=3 --timeout=60 --max-jobs=500 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/behome/worker-emails.log
stopwaitsecs=60

[program:behome-worker-default]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/behome/artisan queue:work database --queue=default --sleep=3 --tries=3 --timeout=60 --max-jobs=500 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/behome/worker-default.log
stopwaitsecs=60
